# deploy/supervisor.conf
[program:put_task_points]
command=php /var/www/demo/bin/console rabbitmq:consumer -m 1000 put_task_points --env=dev -vv
process_name=add_follower_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/var/www/demo/var/log/supervisor.put_task_points.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/var/www/demo/var/log/supervisor.put_task_points.error.log
stderr_capture_maxbytes=1MB
